# pl-files/etc/inittab
::sysinit:/etc/init.d/rcS
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a
::respawn:/bin/agetty 0 tty1
::respawn:/bin/agetty 0 tty2
::respawn:/bin/agetty 0 tty3
::respawn:/bin/agetty 0 tty4
