# pl-files/etc/inittab
::sysinit:/etc/init.d/rcS
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a
tty1::respawn:/bin/agetty 38400 tty1
tty2::respawn:/bin/agetty 38400 tty2
tty3::respawn:/bin/agetty 38400 tty3
tty4::respawn:/bin/agetty 38400 tty4
